# cfg/supervisor/fluunt.conf
[program:fluunt]
command=
command=fluunt-server --loglevel 20 -i 0.0.0.0 -p 8000 -T 15 -M 25 -c 4 -w /mnt/vod/fluunt -C /pvr/ts -L /var/log/fluunt/flunnt-server.log
autostart=true
autorestart=true
stderr_logfile=/var/log/fluunt.err.log
stdout_logfile=/var/log/fluunt.out.log
