# cfg/supervisor/fluunt-image-server.conf
[program:fluunt-image-server]
command=fluunt-image-server --loglevel 30 -i 0.0.0.0 -p 5555 -w /mnt/images/fluunt/ -l /mnt/images/fluunt/local/ -L /var/log/fluunt/fluunt-image-server.log --accesslogfile=/var/log/fluunt/fluunt-image-server_access.log --defaultcachetime=36000
autostart=true
autorestart=true
stderr_logfile=/var/log/fluunt-image-server.err.log
stdout_logfile=/var/log/fluunt-image-server.out.log
